# razer_control_gui/razerdaemon.service
[Unit]
Description=Razer laptop control daemon

[Service]
Type=simple
ExecStart=/usr/share/razercontrol/daemon

[Install]
WantedBy=default.target
